# dnsmasq/dnsmasq.sv.conf
[program:dnsmasq]
command=/bin/bash -c "/home/dhcpd/dnsmasq.sh"
redirect_stderr=true
