# chef-provision.sv.conf
[program:chef-provision]
command=/bin/bash -c "/chef-provisionin-k8s.sh"
redirect_stderr=true