# atftpd/atftpd.sv.conf
[program:atftpd]
command=/bin/bash -c "/run.sh"
redirect_stderr=true